# installer/linux/deb/worker/DEBIAN/postinst
#!/bin/bash
set -e

# Abilita e avvia il servizio worker
systemctl enable worker.service
systemctl start worker.service
