# installer/linux/deb/worker/DEBIAN/prerm
#!/bin/bash
set -e

# Ferma e disabilita il servizio worker
systemctl stop worker.service
systemctl disable worker.service
