# installer/linux/deb/worker/etc/systemd/system/worker.service
[Unit]
Description=Worker Service
After=network.target

[Service]
ExecStart=/usr/bin/worker
Restart=always

[Install]
WantedBy=multi-user.target